# config/led-announcer.service
[Unit]
Description=LED Announcer HTTP Service
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
User=loaled
WorkingDirectory=/home/loaled/Desktop/loaled/led-announcer
Environment="PATH=/home/loaled/Desktop/loaled/led-announcer/.venv/bin:/usr/local/bin:/usr/bin"
ExecStartPre=/bin/sleep 5
ExecStart=/home/loaled/Desktop/loaled/led-announcer/.venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target


