# config/led-announcer.service
[Unit]
Description=LED Announcer HTTP Service
After=network.target sound.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/led-announcer
Environment="PATH=/home/admin/led-announcer/.venv/bin:/usr/local/bin:/usr/bin"
ExecStartPre=/home/admin/led-announcer/scripts/fix_audio.sh
ExecStart=/home/admin/led-announcer/.venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target


